# crontab
* * * * * root bash /root/job.sh && echo "run job.sh $(date)" >> /var/log/cron.log 2>&1
